# systemd/drouter.service.in
[Unit]
Description=DRouter Service
After=network.target remote-fs.target nss-lookup.target
Documentation=man:dmap(1)

[Service]
Type=simple
ExecStart=@prefix@/sbin/drouterd
PrivateTmp=true
Restart=always

[Install]
WantedBy=multi-user.target
