# systemd/drouter.service.in
@GENERATED_SCRIPT_FILE_WARNING@

[Unit]
Description=DRouter Service
After=network.target remote-fs.target nss-lookup.target
Documentation=man:drouterd(8) man:dmap(1) man:dlist(1)

[Service]
Type=simple
ExecStart=@prefix@/sbin/drouterd
TimeoutStopSec=10
PrivateTmp=true
Restart=always

[Install]
WantedBy=multi-user.target
