# systemd/drouter.socket
[Socket]
ListenStream=0.0.0.0:23883
ListenStream=0.0.0.0:9500

[Install]
WantedBy=sockets.target
